# templates/home/user/.config/systemd/user/kuhbs-workspace-attention.service
[Unit]
Description=KUHBS i3 workspace attention markers
After=graphical-session.target

[Service]
Type=simple
# KUHBS installs its Python package into /usr/share/kuhbs, not site-packages.
Environment=PYTHONPATH=/usr/share/kuhbs
ExecStart=/usr/bin/python3 /home/user/.config/i3/kuhbs-integration/kuhbs-workspace-attention.py
Restart=on-failure
RestartSec=2

[Install]
WantedBy=default.target
